# sentinelboot.service
[Unit]
Description=Sentinel Boot Protection
After=local-fs.target

[Service]
Type=oneshot
# Explicitly run the Python implementation to avoid
# any stale compiled binaries.
ExecStart=/usr/bin/python3 /usr/local/share/sentinelroot/boot_protect.py

[Install]
WantedBy=multi-user.target

