# sentinelboot.service
[Unit]
Description=Sentinel Boot Protection
After=local-fs.target

[Service]
Type=oneshot
# Explicitly run as root since boot partition backups require
# administrative permissions
User=root
# Explicitly run the Python implementation to avoid
# any stale compiled binaries.
ExecStart=/usr/bin/python3 -m sentinelroot.boot_protect

[Install]
WantedBy=multi-user.target

